# rc.local
#!/bin/bash

wait_for_socket() {
    local socket=$(basename "${1}") dir=$(dirname "${1}")
    [[ -d ${dir} ]] || mkdir "${dir}"
    while [[ ! -S ${1} ]] ; do
        inotifywait -q -t 2 -e create --format %f "${dir}"

        [[ ${?} == 1 ]] && sleep 2
    done
}

exec &> /var/log/rc.local.log

if [[ -d /etc/rc.local.d ]] ; then
    for rc in /etc/rc.local.d/* ; do
        [[ -r ${rc} ]] && source ${rc}
    done
fi

